# control/legacy_workspaces/archive_legacy_ShadowTag-v2/antigravity-arch/shadowtagai.service
[Unit]
Description=SHADOWTAGAI Kernel Chain Service
Documentation=https://github.com/ShadowTag-v2/aiyou-fastapi-services
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shadowtagai
Group=shadowtagai
WorkingDirectory=/opt/shadowtagai
EnvironmentFile=/etc/shadowtagai/config.env

# Main process
ExecStart=/opt/shadowtagai/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000

# Restart behavior
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Allow write to logs
ReadWritePaths=/var/log/shadowtagai

# Resource limits
MemoryMax=2G
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowtagai

[Install]
WantedBy=multi-user.target
