# artifacts/corpus-intel/com.antigravity.corpus-intel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.corpus-intel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/pikeymickey/.gemini/antigravity/Monorepo-Uphillsnowball/core/rag_evolve.py</string>
        <string>--loop</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GEMINI_API_KEY</key>
        <string>REPLACE_WITH_KEY</string>
        <key>GEMINI_MODEL</key>
        <string>gemini-2.0-flash</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/pikeymickey/.gemini/antigravity/Monorepo-Uphillsnowball</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/pikeymickey/.gemini/antigravity/Monorepo-Uphillsnowball/artifacts/corpus-intel/loop.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/pikeymickey/.gemini/antigravity/Monorepo-Uphillsnowball/artifacts/corpus-intel/loop.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
